# tor1/bin/tor1.ini
[group:tornadoes]
programs = tornado-8000,tornado-8001

[program:tornado-8000]
command = python3 /home/zyk/tornado1/tor1/run.py --port=8000
directory = /home/zyk/tornado1/tor1
autorestart = true
redirect_stderr = true
stdout_logfile = /home/zyk/tornado1/tor1/log/tornado.log
loglevel = info

[program:tornado-8001]
command = python3 /home/zyk/tornado1/tor1/run.py --port=8001
directory = /home/zyk/tornado1/tor1
autorestart = true
redirect_stderr = true
stdout_logfile = /home/zyk/tornado1/tor1/log/tornado.log
loglevel = info
